# patches/etc/inittab
# /etc/inittab: init configuration for busybox init.
# Boot-time system configuration/initialization script.
#
::sysinit:/lib/jam/bin/jam-sysinit

# /sbin/getty respawn shell invocations for selected ttys.
tty1::respawn:/sbin/getty -nl /lib/jam/bin/jam-autologin 38400 tty1
#tty2::respawn:/sbin/getty 38400 tty2
#tty3::askfirst:/sbin/getty 38400 tty3
#tty4::askfirst:/sbin/getty 38400 tty4
#tty5::askfirst:/sbin/getty 38400 tty5
#tty6::askfirst:/sbin/getty 38400 tty6

# Stuff to do when restarting the init 
# process, or before rebooting.
::restart:/lib/jam/bin/jam-restart
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/lib/jam/bin/jam-shutdown

